# steamos-extension-loader/usr/lib/systemd/user/steamos-extension-loader-cleaner.service
[Unit]
Description=Cleanup Extension Loader

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-loader-cleaner
